# deploy/systemd/beigebox.service
[Unit]
Description=BeigeBox — LLM Middleware Control Plane
Documentation=https://github.com/ralabarge/beigebox
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
User=beigebox
Group=beigebox
WorkingDirectory=/opt/beigebox

# Main process
ExecStart=/opt/beigebox/.venv/bin/beigebox dial

# Clean restart on failure with a short backoff
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=120s
StartLimitBurst=5

# Environment — override with /etc/beigebox/env or /opt/beigebox/.env
EnvironmentFile=-/etc/beigebox/env
EnvironmentFile=-/opt/beigebox/.env

# Optional: master key for encrypted connection tokens
# Environment="BB_MASTER_KEY=<64-hex-chars>"

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/beigebox/data /opt/beigebox/workspace/out /opt/beigebox/logs
ProtectHome=yes
PrivateDevices=yes
CapabilityBoundingSet=

# Resource limits
LimitNOFILE=65536
MemoryMax=4G

# Stdout / stderr go to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=beigebox

[Install]
WantedBy=multi-user.target
